# config/systemd/agent@.service
[Unit]
Description=Agent for user %i
After=basic.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/home/%i
Environment=MAIL=/home/%i/Maildir
ExecStart=/usr/local/bin/run-agent.sh
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agent-%i

# ── Security hardening (disabled for Docker compatibility) ───
# Many systemd security and resource directives fail silently inside
# Docker containers, causing services to go straight to inactive (dead)
# with no journal output. Common causes:
#   - Cgroup: MemoryMax/CPUQuota require cgroup delegation that Docker
#     does not always provide to the container's systemd.
#   - Mount namespaces: ProtectSystem/PrivateTmp conflict with Docker's
#     overlay2 filesystem and mount setup.
#   - Seccomp: RestrictAddressFamilies/RestrictSUIDSGID can conflict
#     with the container's seccomp/namespace configuration.
#
# The Docker container boundary already provides process and filesystem
# isolation. Agents still run as unprivileged users with no sudo access.
#
# To re-enable outside Docker, uncomment:
# NoNewPrivileges=true
# ProtectSystem=full
# PrivateTmp=true
# CapabilityBoundingSet=
# AmbientCapabilities=
# RestrictSUIDSGID=true
# RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
# ProtectKernelTunables=true
# ProtectKernelModules=true
# ProtectControlGroups=true
# ProtectKernelLogs=true
# RestrictNamespaces=true
# MemoryMax=512M
# CPUQuota=50%

[Install]
WantedBy=multi-user.target
