# config/systemd/agent@.service
[Unit]
Description=Agent for user %i
After=basic.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/home/%i
ExecStart=/usr/local/bin/run-agent.sh
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agent-%i

# ── Security hardening ──────────────────────────────────────
# Prevent privilege escalation
NoNewPrivileges=true

# Make /usr and /boot read-only (use 'full' instead of 'strict' because
# 'strict' requires remounting the root filesystem read-only, which fails
# on Docker's overlay2 storage driver)
ProtectSystem=full

# Private /tmp so agents cannot see each other's temp files
PrivateTmp=true

# Restrict the set of capabilities to empty (agent needs none)
CapabilityBoundingSet=
AmbientCapabilities=

# Prevent gaining capabilities via setuid/setgid binaries
RestrictSUIDSGID=true

# Restrict address families to only what's needed (IPv4, IPv6, Unix sockets)
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# ── Kernel-level protections (disabled for Docker compatibility) ──
# These directives fail silently in Docker containers, especially on
# older kernels (4.x), causing the service to stay inactive (dead)
# with no journal output. The container boundary provides equivalent
# kernel-level isolation, so these are redundant here.
#
# To re-enable outside Docker, uncomment the following:
# ProtectKernelTunables=true
# ProtectKernelModules=true
# ProtectControlGroups=true
# ProtectKernelLogs=true
# RestrictNamespaces=true

# Basic resource constraints (tune as needed)
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
