# config/systemd/agent@.service
[Unit]
Description=Agent for user %i
After=basic.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/home/%i
ExecStart=/usr/local/bin/run-agent.sh
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agent-%i

# ── Security hardening ───────────────────────────────────────────────
# Directives are grouped by Docker compatibility tier. See
# docs/systemd-cgroup-docker-compat.md for the full analysis.

# Tier 1: Safe in Docker (kernel flags, no mount/cgroup interaction)
NoNewPrivileges=true
RestrictSUIDSGID=true

# Tier 2: Requires cgroup write access (enabled by cgroup:host + rw mount
# in docker-compose.yml). These enforce per-agent resource limits.
MemoryMax=512M
CPUQuota=50%

# ── Tier 3 & 4: Disabled for Docker compatibility ───────────────────
# These directives conflict with Docker's overlay2 filesystem, mount
# setup, or privileged container mode. The Docker container boundary
# already provides equivalent isolation. Agents still run as
# unprivileged users with no sudo access.
#
# Mount namespace conflicts (Tier 3):
# ProtectSystem=full
# PrivateTmp=true
# RestrictNamespaces=true
#
# Kernel/capability restrictions — conflict with privileged mode (Tier 4):
# ProtectKernelTunables=true
# ProtectKernelModules=true
# ProtectControlGroups=true
# ProtectKernelLogs=true
# CapabilityBoundingSet=
# AmbientCapabilities=
# RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
