# config/systemd/agent@.service
[Unit]
Description=Agent for user %i
After=multi-user.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/home/%i
ExecStart=/usr/local/bin/run-agent.sh
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agent-%i

# Basic resource constraints (tune as needed)
# MemoryMax=512M
# CPUQuota=50%

[Install]
WantedBy=multi-user.target
