# config/skel/.config/systemd/user/example.service
[Unit]
Description=Example User Service
After=default.target

[Service]
Type=simple
ExecStart=/bin/bash -c 'echo "Example service running for user: $(whoami)" >> /home/user/service.log'
Restart=on-failure

[Install]
WantedBy=default.target
